# lmsd.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--

	Service manifest for Lyrion Music Server

-->

<service_bundle type='manifest' name='lmsd'>
	<service name='application/lmsd' type='service' version='1'>
	        <create_default_instance enabled='false' />
		<single_instance/>
       
		<dependency name='name-services' grouping='require_all' restart_on='none' type='service'>
			<service_fmri value='svc:/milestone/name-services' />
		</dependency>
		<dependency name='network' grouping='require_all' restart_on='none' type='service'>
			<service_fmri value='svc:/milestone/network' />
		</dependency>
		<dependency name='local-filesystems' grouping='require_all' type='service' restart_on='none'>
			<service_fmri value='svc:/system/filesystem/local' />
		</dependency>
<!--		<dependency name='mysql' grouping='require_all' type='service' restart_on='none'>
                        <service_fmri value='svc:/application/database/mysql:version_51' />
                </dependency>
-->
		<method_context>
                        <method_credential user='lmsd' group='lmsd' />
                </method_context>   
		
     		<exec_method type='method' name='start' exec='/lib/svc/method/svc-lmsd start' timeout_seconds='60' />
        	<exec_method type='method' name='stop' exec='/lib/svc/method/svc-lmsd stop' timeout_seconds='60' />
		<exec_method type='method' name='restart' exec='/lib/svc/method/svc-lmsd restart' timeout_seconds='60' />
		<exec_method type='method' name='refresh' exec='/lib/svc/method/svc-lmsd refresh' timeout_seconds='60' />

		<property_group name='lmsd' type='application'>
                        <propval name='bindir' type='astring' value='/opt/lyrionmusicserver' /> 
			<propval name='perlbin' type='astring' value='/usr/bin/perl' />
                        <propval name='cachedir' type='astring' value='/var/lib/squeezeboxserver/cache' />
			<propval name='prefsdir' type='astring' value='/var/lib/squeezeboxserver/prefs' />
			<propval name='logdir' type='astring' value='/var/lib/squeezeboxserver/logs' />
                        <propval name='transcoding' type='boolean' value='true' />
                        <propval name='sb1slimp3sync' type='boolean' value='true' />
                        <propval name='infolog' type='boolean' value='true' />
                        <propval name='debuglog' type='boolean' value='false' />
                        <propval name='webgui' type='boolean' value='true' />
			<propval name='perfmon' type='boolean' value='false' />
			<propval name='diag' type='boolean' value='false' />
			<propval name='statistics' type='boolean' value='true' />
			<propval name='httpaddr' type='astring' value='0.0.0.0' />
			<propval name='httpport' type='astring' value='9000' />
			<propval name='cliaddr' type='astring' value='0.0.0.0' />
			<propval name='cliport' type='astring' value='9090' />
			<propval name='playeraddr' type='astring' value='0.0.0.0' />
			<propval name='streamaddr' type='astring' value='0.0.0.0' />
			<propval name='dbtype' type='astring' value='SQLite' />
			<propval name='lctype' type='astring' value='de_DE.UTF-8' />
                        <propval name='charset' type='astring' value='UTF-8' />
			<propval name='value_authorization' type='astring' value='solaris.smf.value.lmsd' />
                        <propval name='action_authorization' type='astring' value='solaris.smf.manage.lmsd' />
                        <!-- sub-process core dumps shouldn't restart session -->
                        <propval name='ignore_error' type='astring' value='core,signal' /> 
                </property_group>
		
		<template>
              		<common_name>
				<loctext xml:lang='C'>Lyrion Music Server</loctext>
               		</common_name>
                	<documentation>
                       		<doc_link name='Lyrion Music Server' uri='https://lyrion.org/'/>
               		</documentation>
        	</template>
	</service>
</service_bundle>
